# net-analyzer/datadog-agent/files/init.d/datadog-dogstatsd
#!/sbin/openrc-run
description="DogStatsd is a small server that aggregates your custom app metrics" # http://docs.datadoghq.com/guides/dogstatsd/
command=/opt/datadog-agent/agent/dogstatsd.py
command_args="--use-local-forwarder"
pidfile=/tmp/dogstatsd.pid # TODO
start_stop_daemon_args="--chdir /opt/datadog-agent --background"

depend() {
	use datadog-ddagent
}
